# launchd/com.user.backlog-board.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.backlog-board.server</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PROJECT__/bin/backlog-board-server</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>BACKLOG_API_KEY</key>
        <string>__BACKLOG_API_KEY__</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/backlog-board/server.out.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/backlog-board/server.err.log</string>
</dict>
</plist>
